# vault-warden.service
[Unit]
Description=Vault Warden - Discord Audit Monitor
After=network.target vault.service

[Service]
Type=simple
# Call the 'audit' subcommand
ExecStart=/usr/local/bin/vault-warden -config /etc/vault-warden.yaml audit
Restart=always
RestartSec=10
User=root

[Install]
WantedBy=multi-user.target
